# systemd/ipa-healthcheck-export.timer
[Unit]
Description=Run ipa-healthcheck export every 5 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=15s
Persistent=true

[Install]
WantedBy=timers.target
